# supervisor/aokb-schedule.conf
[program:aokb-schedule-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/apps/schedule-aokb-app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/apps/schedule-aokb-app/storage/logs/worker.log
stopwaitsecs=3600
